# deploy/launchd/com.zippergen.workflow.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zippergen.workflow</string>

  <key>WorkingDirectory</key>
  <string>/Users/YOU/path/to/zippergen</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>uv</string>
    <string>run</string>
    <string>zippergen</string>
    <string>run</string>
    <string>examples/command_center.py:command_center</string>
    <string>--store</string>
    <string>/Users/YOU/.zippergen/runs/command-center.sqlite</string>
    <string>--llm</string>
    <string>openai:gpt-4o</string>
    <string>--services</string>
    <string>live</string>
    <string>--llm-idle-timeout</string>
    <string>300</string>
    <string>--timeout</string>
    <string>0</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>OPENAI_API_KEY</key>
    <string>REPLACE_ME</string>
    <key>ZIPPERGEN_TELEGRAM_TOKEN</key>
    <string>REPLACE_ME</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/YOU/.zippergen/logs/workflow.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/YOU/.zippergen/logs/workflow.err.log</string>
</dict>
</plist>
